# Backup/ex2/crontab
50 22 * * * /usr/local/bin/backup_home.sh
